# build/systemd/atlasclaw-logrotate.service
[Unit]
Description=Rotate AtlasClaw VM file logs
ConditionPathExists=/etc/logrotate.d/atlasclaw

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate /etc/logrotate.d/atlasclaw
